# project/mysql/templates/xtrabackup.service.j2
[Unit]
Description=xtrabackup service
[Service]
Type=oneshot

# Репозиторий
Environment=REPO={{mysql_local_bak}}

