# project/mysql/templates/xtrabackup.timer.j2
[Unit]
Description=xtrabackup

[Timer]
OnUnitActiveSec={{backup_timer_mysql}}

[Install]
WantedBy=timers.target