# project/backup/templates/borg-backup-check.timer.j2
[Unit]
Description=Borg Backup check

[Timer]
OnUnitActiveSec={{check_timer}}

[Install]
WantedBy=timers.target