# project/backup/templates/borg-backup.timer.j2
[Unit]
Description=Borg Backup

[Timer]
OnUnitActiveSec={{backup_timer}}

[Install]
WantedBy=timers.target